# systemd/github-runner-check.service
[Unit]
Description=Verifica status do GitHub Actions Runner periodicamente

[Service]
Type=oneshot
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/check_github_runner.sh

[Install]
WantedBy=multi-user.target
