# systemd/github-runner-start.service
[Unit]
Description=Start GitHub Actions Runner on boot
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/start_github_runner.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
