# systemd/github-runner-check.timer
[Unit]
Description=Checker do status do GitHub Actions Runner

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
