# systemd/mt5-daily-report.service
[Unit]
Description=Relatório diário unificado do projeto MT5

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/MT5-Process-Core-full
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/daily_report.sh

[Install]
WantedBy=multi-user.target
