# systemd/mt5-backup.service
[Unit]
Description=MT5 Trading DB Backup Service
Documentation=file:///home/felipe/mt5-trading-db/docs/backup.md
After=network-online.target postgresql.service docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/felipe/mt5-trading-db

# Carregar variáveis de ambiente
EnvironmentFile=/etc/default/mt5-backup

# Executar script de backup
ExecStart=/home/felipe/mt5-trading-db/scripts/backup.sh

# Timeout de 30 minutos para backup grande
TimeoutStartSec=1800

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-backup

# Restart policy
Restart=no

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/backups/mt5

[Install]
WantedBy=multi-user.target
