# systemd/mt5-maintenance.timer
[Unit]
Description=MT5 Trading DB - Scheduled Maintenance

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
