# systemd/mt5-update.timer
[Unit]
Description=Agendamento diário da atualização automática do stack MT5

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
