# systemd/mt5-daily-report.service
[Unit]
Description=Relatório diário unificado do projeto MT5

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/mt5-trading-db
ExecStart=/bin/bash /home/felipe/mt5-trading-db/scripts/daily_report.sh

[Install]
WantedBy=multi-user.target
