# systemd/mt5-remote-backup.service
[Unit]
Description=Backup diário completo do ambiente MT5 para servidor remoto
After=mt5-update.service mt5-tests.service

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/MT5-Process-Core-full
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/remote_backup.sh

[Install]
WantedBy=multi-user.target
