# systemd/mt5-update.service
[Unit]
Description=Atualização automática do stack MT5
After=network.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/MT5-Process-Core-full
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/update_stack.sh

[Install]
WantedBy=multi-user.target
