# systemd/mt5-healthcheck.service
[Unit]
Description=MT5 docker healthcheck runner

[Service]
Type=oneshot
ExecStart=/home/felipe/MT5-Process-Core-full/scripts/health_unhealthy_check.sh
StandardOutput=journal
StandardError=journal
