# systemd/mt5-maintenance.service
[Unit]
Description=MT5 Trading DB - Maintenance (compose up + health checks)
Wants=docker.service network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/MT5-Process-Core-full
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/maintenance.sh full
TimeoutStartSec=300
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
