# systemd/mt5-backup.timer
[Unit]
Description=MT5 Trading DB Backup Timer
Documentation=file:///home/felipe/MT5-Process-Core-full/docs/backup.md
Requires=mt5-backup.service


[Timer]
# Executar backup todo sábado às 05:00 AM
OnCalendar=Sat *-*-* 05:00:00

# Se o sistema estiver desligado no horário agendado,
# executar assim que o sistema ligar (até 1 hora de atraso)
Persistent=true

# Offset aleatório de até 15 minutos para evitar picos
RandomizedDelaySec=15min

# Precisão de 1 minuto
AccuracySec=1min

[Install]
WantedBy=timers.target
