# systemd/mt5-backup-api.service
[Unit]
Description=MT5 Backup API (FastAPI + Uvicorn)
After=default.target network.target

[Service]
Type=simple
WorkingDirectory=/home/felipe/MT5-Process-Core-full
Environment=API_PORT=9101
Environment=LOG_DIR=/home/felipe/MT5-Process-Core-full/logs/api
ExecStartPre=/usr/bin/mkdir -p /home/felipe/MT5-Process-Core-full/logs/api
ExecStart=/home/felipe/MT5-Process-Core-full/.venv/bin/uvicorn api.app.main:app --host 0.0.0.0 --port ${API_PORT} --log-level info --proxy-headers
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
