# systemd/mt5-remote-backup.service
[Unit]
Description=Backup diário completo do ambiente MT5 para servidor remoto
After=mt5-update.service mt5-tests.service

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/mt5-trading-db
ExecStart=/bin/bash /home/felipe/mt5-trading-db/scripts/remote_backup.sh

[Install]
WantedBy=multi-user.target
