# systemd/mt5-vuln-check.service
[Unit]
Description=Verificação diária de vulnerabilidades em dependências

[Service]
Type=oneshot
WorkingDirectory=/home/felipe/MT5-Process-Core-full
ExecStart=/bin/bash /home/felipe/MT5-Process-Core-full/scripts/check_vulnerabilities.sh

[Install]
WantedBy=multi-user.target
